# cron.daily/sendemail
#!/bin/sh

cd /opt/lemt/defi/ && sudo ./coinpolice.py

